# packages/serverctl/serverctl-0.2.8.tar.gz/serverctl-0.2.8/serverctl@.service
[Unit]
After=network.target

[Service]
Type=simple
Restart=on-failure
WorkingDirectory=%h
ExecStart=/bin/bash -l -c "serverctl run %i"
StandardOutput=file:%h/.local/share/serverctl/logs/%i.log

[Install]
WantedBy=default.target